# systemd/nfs-server.service.d/ftd-recovery.conf
[Unit]
After=srv-clonezilla\x2dimages.mount
Wants=srv-clonezilla\x2dimages.mount
